# setup/crontab
*/1 * * * * chmod +x /WA-Mod-Pack/*.sh; /WA-Mod-Pack/pull.sh